# private_Library/LaunchAgents/com.user.bluetoothconnect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.bluetoothconnect</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>~/.caleb/bluetooth-connect.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/tmp/.com.apple.user.awake</string>
    </array>
</dict>
</plist>
